# timers/apache/apache_restart.service
[Unit]
Description=restart apache server

[Service]
Type=simple
ExecStart=systemctl restart apache2.service

[Install]
WantedBy=timers.target
