# timers/plausible/plausible.timer
[Unit]
Description=start plausible on boot

[Timer]
OnBootSec=15sec
Persistent=true
Unit=plausible.service

[Install]
WantedBy=timers.target
