# timers/generate/generate.service
[Unit]
Description=update site

[Service]
Type=simple
ExecStart=/home/eprints/bin/generate.sh

[Install]
WantedBy=timers.target
