# timers/backup/backup.service
[Unit]
Description=backup bnarchives

[Service]
Type=simple
ExecStart=/home/setup/03_maintainance/backup/backup.sh

[Install]
WantedBy=timers.target
