# timers/indexer/indexer.service
[Unit]
Description=start eprints indexer

[Service]
Type=simple
ExecStart=/home/eprints/bin/indexer.sh

[Install]
WantedBy=timers.target
