# timers/indexer/indexer.timer
[Unit]
Description=start eprints indexer

[Timer]
OnBootSec=15min
Persistent=true
Unit=indexer.service

[Install]
WantedBy=timers.target
